# system_files/usr/lib/systemd/user/bazzite-user-bootstrap.service
[Unit]
Description=Initial image bootstrap

[Service]
Type=oneshot
ExecStart=/usr/local/libexec/bazzite-user-bootstrap

