# system_files/usr/lib/systemd/user/spacium.service
[Unit]
Description=Various vanilla configurations for Spacium

[Service]
Type=oneshot
ExecStart=/usr/local/libexec/spacium

